# 33-dynamic-web /provisioning/golang/goweb.service
[Unit]
Description=goweb

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/opt/goweb/main

[Install]
WantedBy=multi-user.target
