# 08-systemd/monlog/monlog.timer
[Unit]
Description=Execute monlog.service every 30 sec

[Timer]
OnCalendar=*:*:0/30

[Install]
WantedBy=multi-user.target
